# data/siteguard-check.service
[Unit]
Description=SiteGuard periodic block check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/siteguard --check
