# data/siteguard-check.timer
[Unit]
Description=SiteGuard daily block check timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
